# files/restic-backup@.timer
[Unit]
Description=restic backup to %I

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
